# systemd-units/renew-certbot-certs.service
[Unit]
Description=Renews Let's Encrypt certificates via certbot

[Service]
Type=simple
ExecStart=/usr/bin/certbot renew --quiet
